# supervisor/mysql.conf
[program:mariadb]
command=/usr/sbin/mysqld
autostart=true
autorestart=true
user=root
