# supervisor/nodejs_app.conf
[program:nodeapp]
command=npm start --prefix /app/
autostart=true
autorestart=true
